# backends/debian/packages/k3s-system-config/debian/disable-swap.service
[Unit]
Description=Disable swap (required for Kubernetes)
DefaultDependencies=no
Before=swap.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/sbin/swapoff -a
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
